# dietpi/etc/systemd/system/unbound.service.d/cache.conf
[Service]
ExecStop=/etc/unbound/cache-save.sh
